# supervisord.conf
[unix_http_server]
file=/tmp/winefornia-supervisor.sock

[supervisord]
logfile=%(here)s/logs/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/tmp/winefornia-supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/winefornia-supervisor.sock

[program:mail-watcher]
command=%(here)s/.venv/bin/python3 scripts/tastingroom_mail_watcher.py
directory=%(here)s
autostart=true
autorestart=true
startretries=20
startsecs=5
stopwaitsecs=10
stdout_logfile=%(here)s/logs/mail-watcher.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(here)s/logs/mail-watcher-err.log
stderr_logfile_maxbytes=5MB

[program:web]
command=%(here)s/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
directory=%(here)s
autostart=true
autorestart=true
startretries=20
startsecs=5
stopwaitsecs=10
stdout_logfile=%(here)s/logs/web.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(here)s/logs/web-err.log
stderr_logfile_maxbytes=5MB
